# etc/services.d/caddy/run
#!/usr/bin/execlineb -P
with-contenv
s6-setuidgid caddy
/usr/bin/caddy run -config /etc/Caddyfile --adapter caddyfile